# debian/grub-btrfs.path
[Unit]
Description=Monitors for new Timeshift snapshots

[Path]
PathModified=/run/timeshift

[Install]
WantedBy=multi-user.target
